# files/etc/services.d/apache2/run
#!/usr/bin/execlineb -P
/usr/sbin/apachectl -D FOREGROUND
