# files/etc/services.d/pure-ftpd/run
#!/usr/bin/execlineb -P

pure-ftpd /etc/pure-ftpd.conf
